# service/org.nomad.firewall.service
[Unit]
Description=Maui Threes Steps Firewall

[Service]
ExecStart=nomad-service

[Install]
WantedBy=multi-user.target